# coturn/rootfs/etc/services.d/coturn/finish
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: coturn
# Take down the S6 supervision tree when coturn fails
# ==============================================================================
if [[ "${1}" -ne 0 ]] && [[ "${1}" -ne 256 ]]; then
  bashio::log.warning "coturn crashed, halting add-on"
  /run/s6/basedir/bin/halt
fi

bashio::log.info "coturn stopped, restarting..."
